# gif.service.example
# Example systemd unit for the GIF gallery app.
#
# Assumes:
#   - App deployed to /srv/gif/ with a venv at /srv/gif/.venv/
#   - A dedicated 'gif' user/group owns /srv/gif/
#   - uvicorn installed in the venv (uv add uvicorn[standard])
#   - Socket path matches nginx.conf.example upstream
#
# To use:
#   1. Copy to /etc/systemd/system/gif.service
#   2. Adjust paths, user, and worker count as needed
#   3. systemctl daemon-reload
#   4. systemctl enable --now gif

[Unit]
Description=GIF Gallery (uvicorn)
After=network.target

[Service]
Type=simple
User=gif
Group=gif
WorkingDirectory=/srv/gif/gif
RuntimeDirectory=gif

ExecStart=/srv/gif/.venv/bin/uvicorn gif.asgi:application \
    --uds /run/gif/uvicorn.sock \
    --workers 2 \
    --access-log \
    --loop uvloop \
    --http httptools

KillMode=mixed
Restart=on-failure
RestartSec=5

# Hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/srv/gif/gif/media /srv/gif/gif/db.sqlite3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
